# install/retropie-ha.service
[Unit]
Description=RetroPie Home Assistant Integration
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/bin/python3 /home/pi/.config/retropie-ha/status_reporter.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target